# .deployment/x0xd.service
[Unit]
Description=x0xd — local x0x agent daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=exec

# Binary location (testnet/server deployment)
ExecStart=/opt/x0x/x0xd --config /etc/x0x/x0xd.toml

# Working directory
WorkingDirectory=/opt/x0x

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=x0xd

# Restart policy
Restart=always
RestartSec=5s

# Environment
Environment=RUST_LOG=info

# Resource limits
LimitNOFILE=4096

[Install]
WantedBy=default.target
